# deploy/systemd/dimensio-api.service
[Unit]
Description=Dimensio API Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/root/workspace/dimensio
Environment="PATH=/root/workspace/dimensio/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_APP=api/server.py"
Environment="FLASK_ENV=production"
ExecStart=/root/workspace/dimensio/venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 4 \
    --threads 2 \
    --timeout 600 \
    --access-logfile /var/log/dimensio/access.log \
    --error-logfile /var/log/dimensio/error.log \
    --log-level info \
    api.server:app

# 重启策略
Restart=always
RestartSec=10

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

# 资源限制
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
